# scripts/evopaimo-relay.service
[Unit]
Description=EvoPaimo Relay Connector
Documentation=https://github.com/Neon-Wang/xiachong-relay-connect
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/.openclaw/workspace
EnvironmentFile=%h/.config/evopaimo/relay.env
Environment="OPENCLAW_CLI=openclaw"
Environment="OPENCLAW_SESSION_LABEL=xc-relay"
ExecStart=/usr/bin/python3 -u %h/.openclaw/workspace/evopaimo-connect.py --relay ${EVOPAIMO_RELAY_URL} --link-code ${EVOPAIMO_LINK_CODE} --secret ${EVOPAIMO_SECRET} --label ${OPENCLAW_SESSION_LABEL} --agent-file %h/.config/evopaimo/agent.json
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
